# docker/rack-agent.target
[Unit]
Description=Rack Agent Boot Target
Requires=basic.target
After=basic.target
AllowIsolate=yes
